# ssh-agent/.config/systemd/user/ssh-agent.service.d/override.conf
[Service]
Environment=SSH_ASKPASS_REQUIRE=force
Environment=SSH_ASKPASS=%h/.local/bin/ssh-agent-confirm.sh
